# cleanup.timer
[Unit]
Description=Run cleanup.service periodically

[Timer]
OnBootSec=5min
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target