# cleanupd-indexer.service
[Unit]
Description=File Indexing Service for Cleanup (Stage 1: Scan Disk)
After=network.target

[Service]
Type=simple
# Run as a non-privileged user
User=cleanupd
Group=cleanupd

Environment=VIRTUAL_ENV=/opt/cleanup/venv
Environment=PATH=/opt/cleanup/venv/bin:/usr/bin:/bin
ExecStart=/opt/cleanup/venv/bin/python /opt/cleanup/indexer.py --config /opt/cleanup/config.yaml
Restart=on-failure
WorkingDirectory=/opt/cleanup
# Creates /run/cleanupd for pid files, owned by 'cleanupd'
RuntimeDirectory=cleanupd

[Install]
WantedBy=multi-user.target